# 昆虫识别检测模型/deploy/app/supervisord.conf
[supervisord]
nodaemon=true
[program:s1]
command=hub serving start -m kunchongjiance -p 9001
[program:s2]
command=python3 /app/httpServer.py
